# etc/profile.d/a2aliases.sh
alias a2l='a2list'
alias a2mk='a2mksite'
alias a2rl='a2reload'
alias a2rm='a2rmsite'
alias a2rs='a2restart'
alias a2s='a2start'
alias a2sp='a2stop'
alias a2st='a2status'
