# etc/profile.d/a2path.sh
export A2PATH="/etc/apache2"
export A2SITES="$A2PATH/sites-available"
export A2LOG="/var/log/apache2"
export A2ERRORLOG="$A2LOG/error.log"
export A2ACCESSLOG="$A2LOG/access.log"
